# ops/systemd/user/c2-execution-observer.service
[Unit]
Description=C2 IB Execution Observer v1 (PAPER) - broker_event_log.v1.jsonl writer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/node/constellation_2_runtime

# Fail-closed: if python/ibapi not available, service will fail and restart; supervisor/gates will catch missing broker logs.
ExecStart=/usr/bin/env python3 /home/node/constellation_2_runtime/ops/ib/c2_execution_observer_v1.py --host 127.0.0.1 --port 4002 --client-id 79 --poll-seconds 10 --environment PAPER --log-root constellation_2/runtime/truth/execution_evidence_v1/broker_events

Restart=always
RestartSec=3

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
